# deploy/systemd/at-yaris-live-results.timer
[Unit]
Description=Check All Program Track Results Every Five Minutes

[Timer]
OnCalendar=*-*-* 09..23:00/5:00 Europe/Istanbul
Persistent=true
Unit=at-yaris-live-results.service

[Install]
WantedBy=timers.target
